# dist/upstart/solarsan-cli.conf
# solarsan-cli -- SolarSan CLI Daemon
#
# ~trevorj 02/20/13

description	"SolarSan CLI Daemon"

start on solarsan-start
stop on solarsan-stop

respawn
respawn limit 10 5

console log

env HOME=/root

chdir /opt/solarsan
exec ./bin/env ./bin/rpyc_cli
