# dist/upstart/solarsan-storage2.conf
# solarsan-rpc -- SolarSan Storage2 Daemon
#
# ~trevorj 09/01/12

description	"SolarSan Storage2 Daemon"

start on solarsan-start
stop on solarsan-stop

respawn
respawn limit 10 5

console log

env HOME=/root

chdir /opt/solarsan
exec ./bin/env ./bin/rpyc_server2
