# dist/upstart/solarsan-dkv.conf
# solarsan-dkv -- SolarSan Distributed Key-Value Store Daemon
#
# ~trevorj 04/21/13

description	"SolarSan Distributed Key-Value Store Daemon"

start on solarsan-start
stop on solarsan-stop

respawn
respawn limit 10 5

console log

env HOME=/root

chdir /opt/solarsan
exec ./bin/env-reloader ./bin/dkv_serv
